# myapps/install/Madsonic/madsonic.service
[Unit]
Description=Madsonic Media Server
After=network.target

[Service]
User=Madsonic
Group=Madsonic
Type=forking
PreExecStart=/usr/share/madsonic/madsonic.sh

Restart=on-failure
GuessMainPID=no

[Install]
WantedBy=multi-user.target
