# myapps/install/Madsonic/madsonic.service
[Unit]
Description=Madsonic Media Server

[Service]
User=Madsonic
Group=Madsonic
Type=simple
ExecStart=/usr/bin/madsonic --host=0.0.0.0 --port=4040
PIDFile=/run/madsonic/madsonic.pid

[Install]
WantedBy=multi-user.target
