# myapps/install/System/Installserver.service
[Unit]
Wants=network-online.target
After=network-online.target

[Service]
User=root
Type=oneshot
TTYPath=/dev/tty2
ExecStart= /opt/FirstbootInstall.sh

[Install]
WantedBy=multi-user.target
